# mcp-server.service
[Unit]
Description=AI Team Hub MCP Server (SSE mode)
After=network.target

[Service]
Type=simple
User=liunx
WorkingDirectory=/home/liunx/workspace/ai-team-hub
ExecStart=/home/liunx/workspace/ai-team-hub/backend/venv/bin/python3 /home/liunx/workspace/ai-team-hub/mcp-server.py --sse --port 8100
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
